# debian/network-manager.postinst
#!/bin/sh
set -e

case "$1" in
  configure)
  	if [ -x "/etc/init.d/dbus" ]; then
		invoke-rc.d dbus force-reload || true
	fi
						 
	 ;;
esac

#DEBHELPER#
